# home/dot_config/autostart/discord.desktop
[Desktop Entry]
Exec=/usr/bin/discord --start-minimized
Name=discord
Type=Application
