# home/dot_config/autostart/pamac-tray.desktop
[Desktop Entry]
Exec=/usr/bin/pamac-tray
Name=pamac-tray
Type=Application
