# home/dot_config/autostart/steam.desktop
[Desktop Entry]
Exec=/usr/bin/steam
Name=steam
Comment=a launcher for the Steam client from Valve Corporation
Type=Application
