# home/dot_config/systemd/user/autotiling.service
[Unit]
Description=autotiling
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/autotiling-rs

[Install]
WantedBy=sway-session.target
