# home/dot_config/systemd/user/autotiling.service
[Unit]
Description=autotiling
[Service]
ExecStart=/usr/bin/autotiling
