# home/dot_config/systemd/user/dex.service
[Unit]
Description=dex
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/dex -a

[Install]
WantedBy=sway-session.target
