# home/dot_config/systemd/user/dex.service
[Unit]
Description=dex
[Service]
ExecStart=/usr/bin/dex -a
