# home/dot_config/systemd/user/discord.service
[Unit]
Description=discord
[Service]
ExecStart=/usr/bin/discord --start-minimized
