# home/dot_config/systemd/user/dunst.service
[Unit]
Description=A customizable and lightweight notification-daemon
[Service]
ExecStart=/usr/bin/dunst
