# home/dot_config/systemd/user/kitty.service.tmpl
[Unit]
Description=kitty
[Service]
ExecStart=/usr/bin/{{- .term }}
